# systemd/web.service
[Unit]
Description=Web Control Panel
After=mjpg-streamer.service

[Service]
Type=simple
User=tomer
ExecStart=/usr/local/bin/web
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
